# scripts/dbkrab.service
[Unit]
Description=dbkrab - MSSQL CDC Replication Service
After=network.target
Documentation=https://github.com/cnlangzi/dbkrab

[Service]
Type=simple
User=dayi
Group=dayi
WorkingDirectory=/opt/dbkrab
ExecStart=/opt/dbkrab/dbkrab
Environment=HOME=/opt/dbkrab
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
